# rootfs/etc/services.d/openiotai/run
#!/bin/sh
set -x
echo "RUN SCRIPT STARTED" >&2
python3 -u /app/run.py || echo "PYTHON EXIT CODE: $?" >&2
sleep 300
