# rootfs/etc/services.d/openiotai/run
#!/bin/sh
echo "S6 SERVICE STARTED" >&2
exec python3 /app/run.py
